# config/supervisord.ini
[program:rmc-flask]
command=/var/lib/rmc_home/rmc/bin/server_prod.sh
directory=/var/lib/rmc_home/rmc
environment=PATH=/var/lib/rmc_home/rmc/env/bin:%(ENV_PATH)s
priority=11
autostart=true
autorestart=true
killasgroup=true
startretries=2147483648
user=rmc
umask=022
stdout_logfile_maxbytes=4KB
stdout_logfile_backups=0
stderr_logfile_maxbytes=4KB
stderr_logfile_backups=0

[program:rmc-workers]
command=/var/lib/rmc_home/rmc/bin/start_workers.sh
directory=/var/lib/rmc_home/rmc
environment=PATH=/var/lib/rmc_home/rmc/env/bin:%(ENV_PATH)s
priority=11
autostart=true
autorestart=true
killasgroup=true
startretries=2147483648
user=rmc
umask=022
stdout_logfile_maxbytes=4KB
stdout_logfile_backups=0
stderr_logfile_maxbytes=4KB
stderr_logfile_backups=0

[group:rmc]
programs=rmc-workers,rmc-flask
